# docker/supervisor.ini
[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=/srv/script/cron
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
umask=002
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=3600

[program:laravel-web-jobs]
process_name=%(program_name)s_%(process_num)02d
command=php /srv/artisan horizon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
umask=002
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=3600

[program:nginx]
process_name=%(program_name)s_%(process_num)02d
command=/usr/sbin/nginx -g 'daemon off;' -e /dev/stdout
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
umask=002
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=3600
startsecs=1
startretries=1

[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
umask=002
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=3600
startsecs=1
startretries=1
